# rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv execlineb

php-fpm -F -y /usr/local/etc/php-fpm.conf
